# systemd-user/.config/systemd/user/clean-arch.service
[Unit]
Description=Clean Arch caches and temporary files

[Service]
Type=oneshot
ExecStart=%h/dotfiles/scripts/clean_arch.fish
